# bin/bin/other/rfkill_bluetooth.service
[Unit]
Description=example systemd service unit file.

[Service]
ExecStart=rfkill block bluetooth

[Install]
WantedBy=multi-user.target
